# theme.timer
[Unit]
Description = "Run theme switcher"

[Timer]
OnCalendar = daily
Persistent = true

[Install]
WantedBy = timers.target
